# https-dns-proxy.service
[Unit]
Description=Resolve DNS queries over HTTPS. 

[Service]
ExecStart=/usr/local/bin/https_dns_proxy